# docs/deployment/supervisor-queue.conf
# ────────────────────────────────────────────────────────────────────────
# Photo Gallery — Supervisor queue worker config
# ────────────────────────────────────────────────────────────────────────
# Copy to /etc/supervisor/conf.d/photo-gallery.conf then:
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl start photo-gallery-worker:* photo-gallery-downloads:*
#
# Verify:    sudo supervisorctl status
# Restart:   sudo supervisorctl restart photo-gallery-worker:*
# Logs:      tail -f /var/log/photo-gallery-worker.log
# ────────────────────────────────────────────────────────────────────────

# Default queue — notifications, emails, misc short jobs.
# 2 processes is enough for most sites. Scale `numprocs` with traffic.
[program:photo-gallery-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/photo-gallery/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --max-jobs=1000 --timeout=300
directory=/var/www/photo-gallery
autostart=true
autorestart=true
startsecs=3
startretries=3
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/photo-gallery-worker.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stopwaitsecs=3600

# Downloads queue — dedicated workers for the long-running BuildOrderZipJob
# (ZIP large orders → upload to R2 → email signed URL). These have a longer
# --timeout because ZIPing 500 photos can take 10+ minutes. Separate queue
# prevents ZIP jobs from blocking short notification jobs.
[program:photo-gallery-downloads]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/photo-gallery/artisan queue:work redis --queue=downloads --sleep=5 --tries=2 --max-time=3600 --max-jobs=50 --timeout=1800
directory=/var/www/photo-gallery
autostart=true
autorestart=true
startsecs=3
startretries=3
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/photo-gallery-downloads.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stopwaitsecs=3600

# Group so one command can restart all of them:
#   sudo supervisorctl restart photo-gallery:*
[group:photo-gallery]
programs=photo-gallery-worker,photo-gallery-downloads
priority=999
